# affirm.service
[Unit]
Description=affirmations

[Service]
ExecStart="%h/bin/affirm"
